# tangd-update.path.in
[Path]
PathChanged=@jwkdir@
MakeDirectory=true
DirectoryMode=0700

[Install]
WantedBy=multi-user.target
