# tang@.service.in
[Unit]
Description=Tang daemon
ConditionDirectoryNotEmpty=@TANG_DB@

[Service]
ExecStart=@libexecdir@/tang-serve
