# tangd@.service.in
[Unit]
Description=Tang daemon
ConditionDirectoryNotEmpty=@cachedir@

[Service]
StandardInput=socket
StandardOutput=socket
StandardError=journal
ExecStart=@libexecdir@/tangd @cachedir@
