# tangd.socket.in
[Unit]
Description=Tang daemon socket
ConditionDirectoryNotEmpty=@cachedir@

[Socket]
ListenStream=80
Accept=true

[Install]
WantedBy=multi-user.target
