# tangd-update.service.in
[Unit]
Description=Tang Update script

[Service]
Type=oneshot
ExecStart=@libexecdir@/tangd-update @jwkdir@ @cachedir@
