# src/srv/tang-keyd.socket.in
[Unit]
Description=Tang key daemon socket
ConditionDirectoryNotEmpty=@TANG_DB@

[Socket]
ListenDatagram=@TANG_PORT@

[Install]
WantedBy=multi-user.target
