# src/srv/tang-keyd.service.in
[Unit]
Description=Tang key daemon
ConditionDirectoryNotEmpty=@TANG_DB@

[Service]
ExecStart=@libexecdir@/tang-keyd
