# units/tangd@.service.in
[Unit]
Description=Tang Server

[Service]
StandardInput=socket
StandardOutput=socket
StandardError=journal
ExecStart=@libexecdir@/tangd @jwkdir@
