# units/tangd@.service.in
[Unit]
Description=Tang Server
Documentation=man:tang(8)

[Service]
StandardInput=socket
StandardOutput=socket
StandardError=journal
ExecStart=@libexecdir@/tangd @jwkdir@
User=@user@
Group=@group@