# units/tangd.socket
[Unit]
Description=Tang Server socket

[Socket]
ListenStream=80
Accept=true

[Install]
WantedBy=sockets.target
