# units/tangd.socket
[Unit]
Description=Tang Server socket
Documentation=man:tang(8)

[Socket]
ListenStream=80
Accept=true

[Install]
WantedBy=sockets.target
