# units/tangd@.service.in
[Unit]
Description=Tang Server

[Service]
StandardInput=socket
StandardOutput=socket
StandardError=journal
EnvironmentFile=@sysconfdir@/tangd/tangd.conf
ExecStart=@libexecdir@/tangd $TANG_JWKDIR $TANG_AUTHDIR
User=@user@
