# units/tangd.socket.in
[Unit]
Description=Tang Server socket
Requires=tangd-update.path
Requires=tangd-keygen.service

[Socket]
ListenStream=80
Accept=true

[Install]
WantedBy=multi-user.target
