# units/tangd.socket.in
[Unit]
Description=Tang Server socket
Documentation=man:tang(8)

[Socket]
ListenStream=80
Accept=true

ExecStartPre=-/usr/bin/chmod --silent 0440 -- @jwkdir@/*.jwk @jwkdir@/.*.jwk
ExecStartPre=-/usr/bin/chown -R @user@:@group@ @jwkdir@

[Install]
WantedBy=sockets.target
