# src/tang-luks-askpassd.path.in
[Unit]
Description=Tang LUKS askpass daemon path
DefaultDependencies=false

[Path]
DirectoryNotEmpty=/run/systemd/ask-password

[Install]
WantedBy=basic.target
