# src/tangd.socket.in
[Unit]
Description=Tang daemon socket
ConditionDirectoryNotEmpty=@tangdbdir@/keys

[Socket]
ListenStream=80
ListenDatagram=80

[Install]
WantedBy=multi-user.target
