# src/tang-luks-askpassd.service.in
[Unit]
Description=Tang LUKS askpass daemon
DefaultDependencies=false

[Service]
ExecStart=@libexecdir@/tang-luks-askpassd
