# src/tangd.socket.in
[Unit]
Description=Tang daemon socket
ConditionDirectoryNotEmpty=@tangdbdir@/keys

[Socket]
ListenStream=80
Accept=true

[Install]
WantedBy=multi-user.target
